# contrib/init/libercoind.conf
description "Bitcoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env LIBERCOIND_BIN="/usr/bin/libercoind"
env LIBERCOIND_USER="libercoin"
env LIBERCOIND_GROUP="libercoin"
env LIBERCOIND_PIDDIR="/var/run/libercoind"
# upstart can't handle variables constructed with other variables
env LIBERCOIND_PIDFILE="/var/run/libercoind/libercoind.pid"
env LIBERCOIND_CONFIGFILE="/etc/libercoin/libercoin.conf"
env LIBERCOIND_DATADIR="/var/lib/libercoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # libercoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$LIBERCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run libercoind."
        echo "The setting must appear in $LIBERCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"LIBERCOIN Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$LIBERCOIND_PIDDIR"
    chmod 0755 "$LIBERCOIND_PIDDIR"
    chown $LIBERCOIND_USER:$LIBERCOIND_GROUP "$LIBERCOIND_PIDDIR"
    chown $LIBERCOIND_USER:$LIBERCOIND_GROUP "$LIBERCOIND_CONFIGFILE"
    chmod 0660 "$LIBERCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$LIBERCOIND_PIDFILE" \
    --chuid $LIBERCOIND_USER:$LIBERCOIND_GROUP \
    --exec "$LIBERCOIND_BIN" \
    -- \
    -pid="$LIBERCOIND_PIDFILE" \
    -conf="$LIBERCOIND_CONFIGFILE" \
    -datadir="$LIBERCOIND_DATADIR" \
    -disablewallet \
    -daemon

